# docker_files/etc/service/motioneye/run
#!/bin/bash

exec meyectl startserver -c /etc/motioneye/motioneye.conf